# scripts/ldm.crontab
CRON_TZ=UTC
PATH=/opt/miniconda2/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
PYTHONPATH=/opt/rtstats/pylib

@reboot ~/bin/ldmadmin stop; ~/bin/ldmadmin start

10 * * * * cd /opt/rtstats/scripts; python assign_hostname_geom.py

15 * * * * cd /opt/rtstats/scripts; python agg_rtstats.py
